# config/roles/uwsgi/files/uwsgi.conf
# Emperor uWSGI script

description "uWSGI Emperor"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --emperor /etc/uwsgi --uid www-data --gid www-data
